# iCompSystemFiles/etc/profile.d/iComp.sh
alias iCompApi-up='systemctl start iCompApi'
alias iCompApi-down='systemctl stop iCompApi'
alias iCompApi-stat='systemctl status iCompApi'
alias iCompGit='/home/omv1/python/scripts/iComp/gitPull.py'