# iCompSystemFiles/etc/profile.d/iComp.sh
alias iCompApi-up='systemctl start iCompApi'
alias iCompApi-down='systemctl stop iCompApi'
alias iCompApi-stat='systemctl status iCompApi'